# containers/bitcoin/etc-supervisor-bitcoin.conf
[program:bitcoin]
command=/usr/local/bin/bitcoind --daemon
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
testnet=1
